# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=uwsgi --ini /opt/app/app.ini

[program:nginx]
command=nginx -g "daemon off;"